# etc/profile.d/z_container_prompt.sh
source /etc/bashrc
export PS1="(CentOS7.3)$PS1"

